# mariadb/fs/etc/systemd/system/backup.service
[Unit]
Description=auto backup mariadb data

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/bash /opt/backup-database.sh
