# mariadb/fs/etc/systemd/system/backup.timer
[Unit]
Description=auto backup mariadb data

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=10min
OnUnitInactiveSec=1d
OnUnitActiveSec=1d
