# samba/fs/etc/systemd/system/smb.service.d/override.conf
[Service]
ExecStartPost=-/usr/bin/echo "smb service startup complete"
