# samba/fs/etc/systemd/system/smb.service.d/override.conf
[Unit]
Requires=prepare.service
After=prepare.service

[Service]
ExecStartPost=-/usr/bin/echo "smb service startup complete"
StandardOutput=journal+console
StandardError=journal+console
