# qbittorrent/fs/etc/systemd/system/nginx.service.d/after-file.conf
[Unit]
Requires=update-self-ip.service
After=update-self-ip.service
Before=xvnc0.service i3.service
