# qbittorrent/fs/etc/systemd/system/nginx.service.d/create-bridge.conf
[Service]
ExecStartPost=+/usr/bin/env bash /opt/scripts/nginx.start.sh
ExecStopPre=+/usr/bin/env bash /opt/scripts/nginx.stop.sh
