# qbittorrent/fs/etc/systemd/system/nginx.service.d/ensure-oldfile.conf
[Service]
ExecStartPre=+/usr/bin/rm -f /run/nginx/sockets/qbittorrent-admin.sock /run/nginx/sockets/qbittorrent-tracker.sock
