# docker-registry/fs/etc/periodic/daily/gc.sh
#!/bin/sh

registry garbage-collect /etc/docker/registry/config.yml --delete-untagged
