# nextcloud/fs/etc/systemd/system/php-fpm.service.d/override.conf
[Service]
ExecStartPre=!/bin/bash /opt/load-server.sh
EnvironmentFile=/run/.userenvironments
