# gateway-network/fs/usr/lib/systemd/system/ddns.timer
[Unit]
Description=run ddns every 5 min

[Install]
WantedBy=timers.target

[Timer]
OnUnitInactiveSec=5min
OnBootSec=10s
