# virtual-gateway/fs/etc/periodic/15min/run-ddns
#!/bin/sh

set -e

echo "==============================================" >&2
bash -Eeuo pipefail /opt/ddns6.sh
echo "==============================================" >&2
bash -Eeuo pipefail /opt/ddns4.sh
echo "==============================================" >&2
