# nginx/fs/usr/lib/systemd/system/reload.timer
[Unit]
After=nginx.service
Requires=nginx.service

[Install]
WantedBy=multi-user.target

[Timer]
OnUnitInactiveSec=30min
OnBootSec=1s
Unit=reload-handler.service
