# mariadb/fs/etc/systemd/system/mariadb.service.d/socket-path.conf
[Service]
Environment=MYSQLD_OPTS=--skip-name-resolve --socket /run/nginx/sockets/mariadb.sock
