# cloud/fs/etc/crontabs/media_rw
# min   hour    day     month   weekday command
*/5     *       *       *       *       php -d memory_limit=2G /usr/share/webapps/nextcloud/cron.php
