# homeassistant/fs/etc/s6-overlay/s6-rc.d/udhcpc6/up
#!/command/execlineb -P

udhcpc6 -i eth0 -p /tmp/udhcpc6.pid -R -b
