# homeassistant/fs/etc/s6-overlay/s6-rc.d/udhcpc/up
#!/command/execlineb -P

udhcpc -i eth0 -p /tmp/udhcpc.pid -R -b
