# qbittorrent/fs/usr/lib/systemd/system/update-self-ip.timer
[Unit]
Description=update myself ip address to nginx proxy

[Timer]
OnActiveSec=2m

[Install]
WantedBy=timers.target
