# qbittorrent/fs/usr/lib/systemd/system/update-self-ip.timer
[Unit]
Description=update myself ip address to nginx proxy

[Timer]
OnCalendar=*-*-* *:0/5:00

[Install]
WantedBy=timers.target
