# mariadb/fs/etc/systemd/system/nginx-attach.service.d/require.conf
[Unit]
Requires=nginx.service php-fpm.service
After=nginx.service php-fpm.service
